# systemd/bobcat-watchdog.service
[Unit]
Description=BobcatX watchdog checker

[Service]
ExecStart=/etc/monitor-scripts/bobcat-watchdog.sh
