# systemd/pubkeys-check.timer
[Unit]
Description=Scrape the pubkey every hour

[Timer]
OnBootSec=90
OnUnitActiveSec=3600
AccuracySec=1ms

[Install]
WantedBy=timers.target
